# res/win/scheduler-task.xml
<?xml version="1.0" ?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
    <RegistrationInfo>
        <Author>dogkisser</Author>
    </RegistrationInfo>
    <Triggers>
        <LogonTrigger>
            <Enabled>true</Enabled>
            <UserId>{REPLACE_WITH_USERNAME}</UserId>
            <Delay>PT15S</Delay>
        </LogonTrigger>
    </Triggers>
    <Principals>
        <Principal>
            <GroupId>Builtin\Administrators</GroupId>
            <RunLevel>HighestAvailable</RunLevel>
        </Principal>
    </Principals>
    <Settings>
        <Enabled>true</Enabled>
        <AllowStartOnDemand>true</AllowStartOnDemand>
        <AllowHardTerminate>true</AllowHardTerminate>
    </Settings>
    <Actions>
        <Exec>
            <Command>{REPLACE_WITH_GOONTO_BIN}</Command>
            <WorkingDirectory>{REPLACE_WITH_GOONTO_PATH}</WorkingDirectory>
        </Exec>
    </Actions>
</Task>
